# custom-profile.service
[Unit]
Description=Power Service
After=network-online.target

[Service]
ExecStart=/usr/bin/custom-power-manager

[Install]
WantedBy=default.target
